# .apk/artello/gobetween/gobetween.initd
#!/sbin/openrc-run

name=gobetween
supervisor="s6"
s6_service_path="${RC_SVCDIR}/s6-scan/${name}"

depend() {
  need net s6-svscan
  after firewall
}